# devel/etc/catkin/profile.d/10.roslaunch.sh
/home/student/rosws/devel/.private/roslaunch/etc/catkin/profile.d/10.roslaunch.sh